# deploy/systemd/cloudflared-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Flexi ChatGPT App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared --config /etc/cloudflared/flexi-chatgpt-app.yml tunnel run
Restart=always
RestartSec=5
User=cloudflared
Group=cloudflared

[Install]
WantedBy=multi-user.target
